# systemd/my_langgraph.service
[Unit]
Description=LangGraph server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/my_langgraph_agent/agent
ExecStart=/bin/bash -lc 'source %h/my_langgraph_agent/venv/bin/activate && langgraph dev --no-browser'
Environment=PATH=%h/my_langgraph_agent/venv/bin:/usr/local/bin:/usr/bin
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
